# daemon/root/etc/s6-overlay/s6-rc.d/init-mcsm-daemon/run
#!/usr/bin/with-contenv bash

mkdir -p ${MCSM_PATH}/daemon/data
mkdir -p ${MCSM_PATH}/daemon/logs

lsiown -R abc:abc ${MCSM_PATH}/daemon/data
lsiown -R abc:abc ${MCSM_PATH}/daemon/logs
