# web/root/etc/s6-overlay/s6-rc.d/svc-mcsm-web/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
    cd /opt/mcsm/web node --max-old-space-size=8192 --enable-source-maps app.js
